# network-search/ansible_dirs/cron/crontabs/root
# min	hour	day	month	weekday	command
*/5	*	*	*	*	run-parts /data/cron/5min
*/15	*	*	*	*	run-parts /data/cron/15min
*/30	*	*	*	*	run-parts /data/cron/30min
0	*	*	*	*	run-parts /data/cron/hourly
15	6	*	*	*	run-parts /data/cron/daily_morning
15	18	*	*	*	run-parts /data/cron/daily_evening
15	12	*	*	*	run-parts /data/cron/daily_noon
0	2	*	*	*	run-parts /data/cron/daily_overnigt
0	3	*	*	6	run-parts /data/cron/weekly_overnight
0	5	1	*	*	run-parts /data/cron/monthly_overnight

